# docker/etc/service/apache-error-forwarder/run
#!/bin/sh
exec tail -F -n 0 /var/log/apache2/error.log
